# cfg/supervisor/hwi.conf
[program:hwi]
command=python /sprinkle/HWI/hwi.py -h 0.0.0.0
directory=/sprinkle/HWI/
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/sprinkle/HWI/hwi.log
user=root